# kubuntu/theme/theme-setup.timer
[Unit]
Description=Theme Setup Timer

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
