# kubuntu/theme/wallpaper-setup.timer
[Unit]
Description=Wallpaper Setup Timer

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
